# templates/radarr.service.erb
[Unit]
Description=Radarr Daemon
After=network.target

[Service]
User=radarr
Group=radarr

Type=simple
ExecStart=/usr/bin/mono /opt/Radarr/Radarr.exe -nobrowser
TimeoutStopSec=20
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
